# Scripts/pkg-scripts/postinstall
#!/bin/sh
set -eu

DEVSTACK_APP="/Applications/DevStackMenu.app"
IMPORTER_APP="/Applications/Import Compose To DX.app"
DX_BINARY="/usr/local/bin/dx"
LAUNCH_AGENT_FILE="/Library/LaunchAgents/local.devstackmenu.autostart.plist"
LAUNCH_AGENT_LABEL="local.devstackmenu.autostart"
STATE_DIR="/private/tmp/devstackmenu-installer"
POLICY_FILE="$STATE_DIR/policy"
BACKUP_DIR="$STATE_DIR/backup"

cleanup_state() {
  rm -rf "$STATE_DIR"
}

remove_path_if_present() {
  target_path="$1"

  if [ -d "$target_path" ]; then
    rm -rf "$target_path"
  elif [ -e "$target_path" ]; then
    rm -f "$target_path"
  fi
}

restore_or_remove() {
  target_path="$1"
  backup_path="$2"

  remove_path_if_present "$target_path"
  if [ -e "$backup_path" ]; then
    /usr/bin/ditto "$backup_path" "$target_path"
  fi
}

policy="install"
if [ -f "$POLICY_FILE" ]; then
  policy="$(cat "$POLICY_FILE")"
fi

if [ "$policy" = "keep" ]; then
  restore_or_remove "$DEVSTACK_APP" "$BACKUP_DIR/DevStackMenu.app"
  restore_or_remove "$IMPORTER_APP" "$BACKUP_DIR/Import Compose To DX.app"
  restore_or_remove "$DX_BINARY" "$BACKUP_DIR/dx"
  restore_or_remove "$LAUNCH_AGENT_FILE" "$BACKUP_DIR/local.devstackmenu.autostart.plist"
  cleanup_state
  echo "Existing DevStack installation kept unchanged."
  exit 0
fi

cleanup_state

if [ ! -x "$DX_BINARY" ]; then
  echo "Installation failed: /usr/local/bin/dx is not installed" >&2
  exit 1
fi

DX_STATUS_OUTPUT="$(mktemp)"
DX_EXIT_CODE=0

if ! "$DX_BINARY" >"$DX_STATUS_OUTPUT" 2>&1; then
  DX_EXIT_CODE=$?
fi

if [ "$DX_EXIT_CODE" != "0" ] && \
   ! grep -Fq "Usage:" "$DX_STATUS_OUTPUT"; then
  echo "Installation failed: /usr/local/bin/dx is not executable" >&2
  cat "$DX_STATUS_OUTPUT" >&2
  rm -f "$DX_STATUS_OUTPUT"
  exit 1
fi

rm -f "$DX_STATUS_OUTPUT"

if ! "$DX_BINARY" status >/dev/null 2>&1; then
  echo "Installation failed: /usr/local/bin/dx is not executable" >&2
  exit 1
fi

if [ ! -f "$LAUNCH_AGENT_FILE" ]; then
  echo "Installation failed: $LAUNCH_AGENT_FILE is not installed" >&2
  exit 1
fi

chown root:wheel "$LAUNCH_AGENT_FILE" >/dev/null 2>&1 || true
chmod 644 "$LAUNCH_AGENT_FILE" >/dev/null 2>&1 || true

console_user="$(stat -f %Su /dev/console 2>/dev/null || printf '%s' root)"
console_uid=""
if [ -n "$console_user" ] && [ "$console_user" != "root" ] && [ "$console_user" != "loginwindow" ]; then
  console_uid="$(id -u "$console_user" 2>/dev/null || true)"
fi

if [ -n "$console_uid" ]; then
  launchctl bootout "gui/$console_uid/$LAUNCH_AGENT_LABEL" >/dev/null 2>&1 || \
    launchctl bootout "gui/$console_uid" "$LAUNCH_AGENT_FILE" >/dev/null 2>&1 || true

  if launchctl bootstrap "gui/$console_uid" "$LAUNCH_AGENT_FILE" >/dev/null 2>&1; then
    echo "DevStackMenu will launch at login and has been started for $console_user."
  else
    launchctl asuser "$console_uid" /usr/bin/open -gj "$DEVSTACK_APP" >/dev/null 2>&1 || true
    echo "DevStackMenu launch agent installed; immediate start will occur on the next login for $console_user."
  fi
else
  echo "DevStackMenu launch agent installed; immediate start skipped because no console user session was detected."
fi

echo "/usr/local/bin/dx installed and executable."
